# ansible-linux-server-deploy/roles/failover/templates/failover.timer.j2
[Unit]
Description=Periodic automated failover probe

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ failover_timer_interval }}
Persistent=true

[Install]
WantedBy=timers.target
